# buildlinux/debian/prerm
#!/bin/sh

if [ -x /bin/systemctl ]; then
    systemctl disable shieldtra
    systemctl stop shieldtra
fi

#DEBHELPER#
exit 0
