# buildlinux/debian/prerm
#!/bin/sh

if [ -x /bin/systemctl ]; then
    systemctl disable safejumper
    systemctl stop safejumper
fi

#DEBHELPER#
exit 0
